# enable-ir.service
[Unit]
Description=Map IR codes

[Service]
ExecStart=ir-keytable -c -w /home/pi/controller-go/L336

[Install]
WantedBy=default.target
